# deb/DEBIAN/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl stop compose-watcher.service || true
    systemctl disable compose-watcher.service || true
fi
