# deb/DEBIAN/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl stop dispenser.service || true
    systemctl disable dispenser.service || true
fi
